# etc/systemd/system/apache2.service.d/cuems-master-group.conf
[Unit]
PartOf=cuems-master.service

[Install]
WantedBy=cuems-master.service

