# etc/systemd/system/isc-dhcp-server.service.d/cuems-wifi-group.conf
[Unit]
PartOf=cuems-wifi.service

[Install]
WantedBy=cuems-wifi.service

