# systemd/hostapd.service.d/cuems-wifi-group.conf
[Unit]
PartOf=cuems-wifi.service
# Start this unit after the app.service start

[Install]
WantedBy=cuems-wifi.service

