# etc/systemd/system/cuems-nodeconf.service
[Unit]
Description=cuems-nodeconf
After=network.target


[Service]
Type=oneshot
RemainAfterExit=yes
Restart=no
ExecStart=/usr/bin/python3.7 /home/stagelab/src/cuems-engine/src/nodeconf.py
SuccessExitStatus=100 101
TimeoutSec=900

[Install]
WantedBy=default.target